# camilladsp/DEBIAN/postinst
#!/bin/sh
set -e
mkdir -p /var/lib/camilladsp
systemctl daemon-reload
systemctl enable camilladsp
