# camilladsp/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop camilladsp || true
systemctl disable camilladsp || true
