# camillagui/DEBIAN/prerm
#!/bin/sh
set -e
systemctl stop camillagui || true
systemctl disable camillagui || true
